# devops/upstart.conf
description "TaxPlan Web Server"
author      "EllefsonTech"

start on startup
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    exec sudo -u ubuntu start foreman
end script